# sha-bridge/rootfs/etc/services.d/sha_bridge/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

declare serial_port debug mqtt_url mqtt_username mqtt_password

serial_port=$(bashio::config 'serial_port')
if [ "${serial_port}" != "null" ]; then
  export SHA_SERIAL_DEVICE="${serial_port}"
fi

debug=$(bashio::config 'debug')
if [ "${debug}" != "null" ]; then
  export SHA_DEBUG="${debug}"
fi

mqtt_url=$(bashio::config 'mqtt_url')
if [ "${mqtt_url}" != "null" ]; then
  export SHA_MQTT_URL="${mqtt_url}"
fi

mqtt_username=$(bashio::config 'mqtt_username')
if [ "${mqtt_username}" != "null" ]; then
  export SHA_MQTT_USER="${mqtt_username}"
fi

mqtt_password=$(bashio::config 'mqtt_password')
if [ "${mqtt_password}" != "null" ]; then
  export SHA_MQTT_PASSWORD="${mqtt_password}"
fi

echo "*** ENV ***"
env
echo "***********"

exec /usr/bin/sha_bridge
