# sha-bridge/rootfs/etc/services.d/sha_bridge/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

declare program

if bashio::config.has_value 'serial_port'; then
  export SHA_SERIAL_DEVICE="$(bashio::config 'serial_port')"
fi

if bashio::config.true 'debug'; then
  export RUST_LOG="info,sha_bridge=trace"
else
  export RUST_LOG="info"
fi

if bashio::config.true 'mqtt'; then
  declare mqtt_host mqtt_ssl mqtt_port mqtt_scheme
  mqtt_host="$(bashio::services mqtt 'host')"
  mqtt_ssl="$(bashio::services mqtt 'ssl')"
  if [ "${mqtt_ssl}" = "true" ]; then
    mqtt_scheme="mqtts"
  else
    mqtt_scheme="mqtt"
  fi
  mqtt_port="$(bashio::services mqtt 'port')"
  export SHA_MQTT_URL="${mqtt_scheme}://${mqtt_host}:${mqtt_port}?client_id=standaert_ha"
  export SHA_MQTT_USER="$(bashio::services mqtt 'username')"
  export SHA_MQTT_PASSWORD="$(bashio::services mqtt 'password')"
fi

program="$(bashio::config 'program')"
if [ -f "${program}" ]; then
  export SHAL_PROGRAM="${program}"
  if bashio::config.has_value 'upload'; then
    export SHA_UPLOAD="$(bashio::config 'upload')"
  fi
else
  echo "Not loading program, ${program} does not exist"
fi

exec /usr/bin/sha_bridge
